# com.whisperflow.dictation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.whisperflow.dictation</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/charlottezhou/projects/whisper/.venv/bin/python</string>
        <string>/Users/charlottezhou/projects/whisper/main.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/charlottezhou/projects/whisper</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <!-- restart if it crashes, but stay quit after Quit from the menu bar -->
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/whisperflow.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/whisperflow.err</string>
</dict>
</plist>
